# xlink_kai/rootfs/etc/services.d/kai/run
#!/usr/bin/with-contenv bash
echo "[INFO] Avvio di XLink Kai..."
exec /opt/xlinkkai/kaiEngine -d
